# build/linux/debian_sid_i386-sysroot/debian/libnss-db/DEBIAN/postinst
#!/bin/sh

set -e

if [ "$1" = "configure" ]; then
	if [ -x /etc/init.d/nscd ]; then
		invoke-rc.d nscd restart
	fi
fi





exit 0
